# rootfs/etc/s6-overlay/s6-rc.d/postgres/run
#!/command/with-contenv sh
exec sudo -u postgres /usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/data -c config_file=/etc/postgresql/15/main/postgresql.conf